# com.sistema-gerped.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sistema-gerped.backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd "/Users/ericobrandao/Downloads/GERPED PLUS" && /Users/ericobrandao/Downloads/GERPED\ PLUS/scripts/backup_financeiro.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>21600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ericobrandao/Downloads/GERPED PLUS/instance/backups/backup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ericobrandao/Downloads/GERPED PLUS/instance/backups/backup.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
